# templates/launchd-puller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>{{LABEL}}</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{NODE}}</string>
    <string>{{CCM_BIN}}</string>
    <string>pull</string>
  </array>
  <key>StartInterval</key>
  <integer>{{INTERVAL}}</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{REPO}}/.logs/puller.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>{{REPO}}/.logs/puller.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>{{PATH}}</string>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>PARAMEM_REPO</key>
    <string>{{REPO}}</string>
  </dict>
</dict>
</plist>
