# templates/systemd-maintain.service
[Unit]
Description=paramem autonomous maintenance

[Service]
Type=oneshot
Environment=PARAMEM_REPO={{REPO}}
ExecStart={{NODE}} {{CCM_BIN}} maintain
StandardOutput=append:{{REPO}}/.logs/maintain.stdout.log
StandardError=append:{{REPO}}/.logs/maintain.stderr.log
