# templates/systemd-pull.service
[Unit]
Description=paramem pull

[Service]
Type=oneshot
Environment=PARAMEM_REPO={{REPO}}
ExecStart={{NODE}} {{CCM_BIN}} pull
StandardOutput=append:{{REPO}}/.logs/puller.stdout.log
StandardError=append:{{REPO}}/.logs/puller.stderr.log
